# etc/gum.conf
# This is the upstart config I use for https://willnorris.com/
description "Gum Short URL Resolver"
author "Will Norris <will@willnorris.com>"

start on (net-device-up)
stop on runlevel [!2345]

respawn
exec start-stop-daemon --start -c www-data --exec /usr/bin/gum -- \
    -addr localhost:8002 \
    -log_dir /var/log/gum \
    -redirect w=/wiki/ \
    -redirect '+=https://plus.google.com/+willnorris/' \
    -static_dir /var/www/willnorris.com/public
